# examples/com.user.fundtracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.fundtracker</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/FULL/PATH/TO/flash/scripts/fund_tracker.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/FULL/PATH/TO/flash</string>
    
    <!-- Run daily at 8:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/FULL/PATH/TO/flash/logs/fundtracker.log</string>
    
    <key>StandardErrorPath</key>
    <string>/FULL/PATH/TO/flash/logs/fundtracker.error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
